# install/tailscale-eth-tuning.service
[Unit]
Description=Tailscale Ethtool UDP GRO Tuning
After=network.target
# Ensure it waits for the specific device 
Requires=sys-subsystem-net-devices-enp1s0.device
After=sys-subsystem-net-devices-enp1s0.device

[Service]
Type=oneshot
ExecStart=/usr/sbin/ethtool -K enp1s0 rx-udp-gro-forwarding on rx-gro-list off
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
